# root/etc/service/x11vnc/run
#!/bin/bash

/usr/bin/x11vnc -ncache 10 -xkb -forever -noxrecord -noxfixes -noxdamage -display :0
